# etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:rsyslogd]
command = /usr/sbin/rsyslogd -n

[program:postfix]
command=/run.sh
autorestart=true

